# infra/duo/upstart/master/redis-stop.conf
#!upstart
description "Redis stop"

start on stopping redis

task

script
    stop redis-instance PORT=6379
end script