# infra/duo_m/upstart/mongod-start.conf
#!upstart
description "Mongod start"

start on starting mongod

task

script
    sudo -u mongodb mongodb-start.sh
end script
