# infra/duo_m/upstart/nodejs-start.conf
#!upstart
description "Nodejs start"

start on started nodejs

task

script
    sleep 10
    start nodejs-instance PORT=8081
    start nodejs-instance PORT=8082
end script
