# infra/simple/upstart/nodejs-start.conf
#!upstart
description "Nodejs start"

start on starting nodejs

task

script
    start nodejs-instance PORT=8081
end script
