# infra/single/upstart/mongod.conf
# Ubuntu upstart file at /etc/init/mongod.conf

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

