# infra/single/upstart/redis-start.conf
#!upstart
description "Redis start"

start on starting redis

task

script
    start redis-instance PORT=6379
    start redis-instance PORT=6380
    start redis-instance PORT=6381
end script