# infra/single/upstart/redis-stop.conf
#!upstart
description "Redis stop"

start on stopping redis

task

script
    stop redis-instance PORT=6379
    stop redis-instance PORT=6380
    stop redis-instance PORT=6381
end script