# infra/single/upstart/sentinel.conf
#!upstart
description "Sentinel upstart conf"

start on started redis

stop on shutdown
