# infra/duo_m/upstart/slave/sentinel-start.conf
#!upstart
description "Sentinel start"

start on starting sentinel

task

script
    start sentinel-instance PORT=26379
end script