# infra/xuser-start.conf
#!upstart
description "Xuser start"

start on starting xuser

task

script
    start xuser-instance PORT=8081
    start xuser-instance PORT=8082
    start xuser-instance PORT=8083
end script
