# infra/upstart/mongod-stop.conf
#!upstart
description "Mongod stop"

start on stopping mongod

task

script
    stop mongod-instance CONFILE=/data-mongodb/xuser-1/primary.conf PIDFILE=/var/run/mongod_1.pid
    stop mongod-instance CONFILE=/data-mongodb/xuser-2/secondary.conf PIDFILE=/var/run/mongod_2.pid
    stop mongod-instance CONFILE=/data-mongodb/backup/backup.conf PIDFILE=/var/run/mongod_3.pid
end script
