# infra/upstart/mongod-stop.conf
#!upstart
description "Mongod stop"

start on stopping mongod

task

script
    sudo -u mongodb mongodb-stop.sh
end script
