# infra/upstart/nodejs-stop.conf
#!upstart
description "Nodejs stop"

start on stopping nodejs

task

script
    stop nodejs-instance PORT=8081
    stop nodejs-instance PORT=8082
    stop nodejs-instance PORT=8083
end script
