# infra/upstart/nodejs.conf
#!upstart
description "Nodejs upstart conf"

start on (local-filesystems and net-device-up IFACE=eth0 and started mongod)
stop on shutdown
