# infra/upstart/sentinel.conf
#!upstart
description "Sentinel upstart conf"

start on redis
stop on shutdown