# infra/upstart/xuser-stop.conf
#!upstart
description "Xuser stop"

start on stopping xuser

task

script
    stop xuser-instance PORT=8081
    stop xuser-instance PORT=8082
    stop xuser-instance PORT=8083
end script
