# infra/duo/upstart/slave/sentinel-stop.conf
#!upstart
description "Sentinel stop"

start on stopping sentinel

task

script
    stop sentinel-instance PORT=26379
end script