# infra/duo/upstart/nodejs-start.conf
#!upstart
description "Nodejs start"

start on starting nodejs

task

script
    start nodejs-instance PORT=8081
    start nodejs-instance PORT=8082
    start nodejs-instance PORT=8083
end script
